# 24.04/root/etc/skel/.config/autostart/artha.desktop
[Desktop Entry]
Type=Application
Exec=artha &
Hidden=false
Name=Artha
Comment=Offline English dictionary
X-MATE-Autostart-Delay=0
