# 24.04/root/etc/skel/.config/autostart/transmission-gtk.desktop
[Desktop Entry]
Type=Application
Exec=transmission-gtk
Hidden=true
Name=Transmission-gtk
Comment=Download torrent on start
X-MATE-Autostart-Delay=0
